# rust-rebuild/systemd/dev-backup-monthly.service
[Unit]
Description=Monthly dev backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dev-backup ws run-month %i
